# template/repos/disable-repos.cron
# disable all repos if the server is connected to spacewalk. 
30 03 * * * root /usr/local/bin/disable-repos &> /dev/null
